# monitoramento-sistema/monitoramento_sistema.timer
[Unit]
Description=Timer para execução periódica do Monitoramento de Servidor

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
